# service/detect-lights_openwrt.initd
#!/bin/sh /etc/rc.common

# https://wiki.openwrt.org/doc/techref/initscripts

START=99
STOP=01

start() {
	/usr/bin/hueservice-detect-lights start	
}                 
 
stop() {          
	/usr/bin/hueservice-detect-lights stop
}
