# service/openwrt.initd/detect-bridge
#!/bin/sh /etc/rc.common

# https://wiki.openwrt.org/doc/techref/initscripts

START=97
STOP=03

EXTRA_COMMANDS="status"
EXTRA_HELP="	status	Show status informations about the service"

_EXEC=/usr/bin/hueservice-detect-bridge

start() {
	$_EXEC start
}

status() {
	$_EXEC status
}

stop() {
	$_EXEC stop
}
