# service/openwrt.initd/load-default
#!/bin/sh /etc/rc.common

# https://wiki.openwrt.org/doc/techref/initscripts

START=99
STOP=01

start() {
	/usr/bin/hueload-default &
}

stop() {
	/usr/bin/hue kill
}
